# packages/windows_open_package/correlation_rules/mitre_attck_persist/Create_persist_via_Hidden_Run_key_value/tests/raw_events_1.json
{"body":"{\"Event\":{\"xmlns\":\"http://schemas.microsoft.com/win/2004/08/events/event\",\"System\":{\"Provider\":{\"Name\":\"Microsoft-Windows-Sysmon\",\"Guid\":\"{5770385f-c22a-43e0-bf4c-06f5698ffbd9}\"},\"EventID\":\"13\",\"Version\":\"2\",\"Level\":\"4\",\"Task\":\"13\",\"Opcode\":\"0\",\"Keywords\":\"0x8000000000000000\",\"TimeCreated\":{\"SystemTime\":\"2020-07-04T14:18:58.2687126Z\"},\"EventRecordID\":\"306346\",\"Correlation\":\"\",\"Execution\":{\"ProcessID\":\"3400\",\"ThreadID\":\"4136\"},\"Channel\":\"Microsoft-Windows-Sysmon/Operational\",\"Computer\":\"MSEDGEWIN10\",\"Security\":{\"UserID\":\"S-1-5-18\"}},\"EventData\":{\"Data\":[{\"Name\":\"RuleName\",\"text\":\"Persistence - Hidden Run value detected\"},{\"Name\":\"EventType\",\"text\":\"SetValue\"},{\"Name\":\"UtcTime\",\"text\":\"2020-07-04 14:18:58.231\"},{\"Name\":\"ProcessGuid\",\"text\":\"{747f3d96-8fd2-5f00-0000-0010c15d2200}\"},{\"Name\":\"ProcessId\",\"text\":\"3728\"},{\"Name\":\"Image\",\"text\":\"C:\\\\Users\\\\Public\\\\tools\\\\evasion\\\\a.exe\"},{\"Name\":\"TargetObject\",\"text\":\"HKLM\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\\"},{\"Name\":\"Details\",\"text\":\"\\\"c:\\\\windows\\\\tasks\\\\taskhost.exe\\\"\"}]}}}","recv_ipv4":"127.0.0.1","recv_time":"2023-06-12T19:17:17.180Z","task_id":"00000000-0000-0000-0000-000000000000","tag":"some_tag","mime":"application/x-pt-eventlog","normalized":false,"input_id":"00000000-0000-0000-0000-000000000000","type":"raw","uuid":"7e5ba700-f8f5-4874-a4ac-2a91c9dddd47"}
{"body":"{\"Event\":{\"xmlns\":\"http://schemas.microsoft.com/win/2004/08/events/event\",\"System\":{\"Provider\":{\"Name\":\"Microsoft-Windows-Sysmon\",\"Guid\":\"{5770385f-c22a-43e0-bf4c-06f5698ffbd9}\"},\"EventID\":\"13\",\"Version\":\"2\",\"Level\":\"4\",\"Task\":\"13\",\"Opcode\":\"0\",\"Keywords\":\"0x8000000000000000\",\"TimeCreated\":{\"SystemTime\":\"2023-06-12T18:30:51.2794969Z\"},\"EventRecordID\":\"48185854\",\"Correlation\":\"\",\"Execution\":{\"ProcessID\":\"3636\",\"ThreadID\":\"4616\"},\"Channel\":\"Microsoft-Windows-Sysmon/Operational\",\"Computer\":\"win10-work.stand2008.local\",\"Security\":{\"UserID\":\"S-1-5-18\"}},\"EventData\":{\"Data\":[{\"Name\":\"RuleName\",\"text\":\"T1547_001,Create Persistance: Registry Run Keys\"},{\"Name\":\"EventType\",\"text\":\"SetValue\"},{\"Name\":\"UtcTime\",\"text\":\"2023-06-12 18:30:51.264\"},{\"Name\":\"ProcessGuid\",\"text\":\"{2b856446-645b-6487-ba09-00000000c100}\"},{\"Name\":\"ProcessId\",\"text\":\"2460\"},{\"Name\":\"Image\",\"text\":\"C:\\\\Windows\\\\system32\\\\reg.exe\"},{\"Name\":\"TargetObject\",\"text\":\"HKLM\\\\SOFTWARE\\\\WOW6432Node\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\CalcImplant\"},{\"Name\":\"Details\",\"text\":\"C:\\\\Windows\\\\System32\\\\calc.exe\"},{\"Name\":\"User\",\"text\":\"STAND2008\\\\Администратор\"}]}}}","recv_ipv4":"127.0.0.1","recv_time":"2023-06-12T19:17:17.180Z","task_id":"00000000-0000-0000-0000-000000000000","tag":"some_tag","mime":"application/x-pt-eventlog","normalized":false,"input_id":"00000000-0000-0000-0000-000000000000","type":"raw","uuid":"b4349995-7d1f-4292-9177-fe3ed9b7436d"}